# php71-apache/etc/supervisor/conf.d/apache2.conf
[program:apache2]
command=/usr/sbin/apache2
stdout_logfile=/var/log/apache2/apache2.log
redirect_stderr=true
autostart=true
autorestart=true

